# deploy/gittimes-api.service
[Unit]
Description=GitTimes agent API (x402-paid REST surface)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chris
WorkingDirectory=/home/chris/gittimes
EnvironmentFile=/home/chris/gittimes/deploy/api.env
ExecStart=/usr/bin/node api-server.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
